# services.d/github-runner/run
#!/usr/bin/with-contenv bash
exec $RUNNER_PATH/bin/Runner.Listener run \
     --startuptype service